# dockers/docker-gcu/supervisord.conf
[supervisord]
logfile_maxbytes=1MB
logfile_backups=2
nodaemon=true

[eventlistener:dependent-startup]
command=python3 -m supervisord_dependent_startup
autostart=true
autorestart=unexpected
startretries=0
exitcodes=0,3
events=PROCESS_STATE
buffer_size=1024

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n -iNONE
priority=1
autostart=false
autorestart=unexpected
stdout_logfile=NONE
stdout_syslog=true
stderr_logfile=NONE
stderr_syslog=true
dependent_startup=true

# One-shot program: create the GCU venv and install the wheel.
# priority=2 ensures rsyslogd is up first.
# startsecs=0 means supervisord does not wait for it to stay running.
# autorestart=false because it exits intentionally on success.
[program:gcu-setup]
command=/usr/bin/gcu-init.sh setup
priority=2
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=rsyslogd:running

# Health-check program: verifies venv GCU and host GCU respond correctly.
# Waits for gcu-setup to have exited (supervisord_dependent_startup
# 'exited' state).  The docker-init.sh healthcheck subcommand additionally
# checks for the .setup_complete sentinel file so it exits cleanly if
# gcu-setup failed before the file was written.
[program:gcu-healthcheck]
command=/usr/bin/gcu-init.sh healthcheck
priority=3
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=gcu-setup:exited
