# platform/broadcom/sonic-platform-modules-nexthop/common/etc/systemd/system/pmon.service.d/override.conf
[Unit]
After=pddf-platform-init.service
StartLimitIntervalSec=0
StartLimitBurst=0
